# mcp-coordination/deploy/mcp-coordination.service
[Unit]
Description=MCP Coordination Server (Redis-backed distributed locking)
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=cooneycw
WorkingDirectory=/home/cooneycw/Projects/claude-power-pack/mcp-coordination
ExecStart=/home/cooneycw/Projects/claude-power-pack/mcp-coordination/start-server.sh --daemon
Restart=on-failure
RestartSec=5

# Environment - uv should be in PATH
Environment="PATH=/home/cooneycw/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-coordination

[Install]
WantedBy=multi-user.target
