# extras/redis-coordination/mcp-server/deploy/mcp-coordination.service
[Unit]
Description=MCP Coordination Server (Redis-backed distributed locking)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/claude-power-pack/mcp-coordination
ExecStart=%h/.local/bin/uv run python src/server.py
Restart=on-failure
RestartSec=5

# Environment - uv and local bin in PATH
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-coordination

[Install]
WantedBy=default.target
