# extras/redis-coordination/mcp-server/deploy/mcp-coordination.service
# Install as a user service: systemctl --user enable mcp-coordination
# Copy to: ~/.config/systemd/user/mcp-coordination.service
# %h expands to the user's home directory in user services.

[Unit]
Description=MCP Coordination Server (Redis-backed distributed locking)
After=network.target redis-server.service

[Service]
Type=simple
WorkingDirectory=%h/Projects/claude-power-pack/extras/redis-coordination/mcp-server
EnvironmentFile=%h/Projects/claude-power-pack/extras/redis-coordination/mcp-server/.env
ExecStart=%h/.local/bin/uv run python src/server.py
Restart=on-failure
RestartSec=5

# Environment - uv and local bin in PATH
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-coordination

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=default.target
